# projects/water-plant-login/systemd/water-plant-backend.service
[Unit]
Description=Water Plant Backend Server
After=network.target

[Service]
Type=simple
User=hzh
WorkingDirectory=/home/hzh/文档/project/projects/water-plant-login
ExecStart=/usr/bin/node server/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
